# templates/default/mailcatcher.upstart.conf.erb
#!upstart
description "Mailcatcher SMTP dev tool"

start on started

script
  exec /opt/vagrant_ruby/bin/mailcatcher --http-ip="0.0.0.0"
end script
